# systemd/mi-web-chat.service.d/10-mi-runtime.conf
# Managed by Mi stack installer. Do not place secrets in this file.
[Service]
Environment=MI_GATEWAY_CLIENT=/home/kyle/.local/share/mi/mi-gateway-client.py
Environment=PI_CMD=/home/kyle/bin/pi-gateway
Environment=PATH=/home/kyle/.nvm/versions/node/v24.15.0/bin:/usr/local/bin:/usr/bin:/bin
