# systemd/mi-web-chat.service.d/10-nvm-pi-path.conf
# Mi V2 invokes /home/kyle/bin/pi-gateway, which execs `pi` via PATH.
# Keep it on the supported NVM Pi binary rather than the distro-global Pi.
[Service]
Environment=PATH=/home/kyle/.nvm/versions/node/v24.15.0/bin:/usr/local/bin:/usr/bin:/bin
