# gateway/llm-gateway.service.d/20-codex-subscription.conf
# Installed by scripts/install-mi-subscription-gateway-root.sh.
# Pi's subscription storage belongs to kyle; only its agent directory is writable.
[Service]
User=kyle
Group=kyle
Environment=HOME=/home/kyle
LoadCredential=
LoadCredential=gateway:/etc/agent-secrets/local-agent-gateway.token
ProtectHome=read-only
ReadWritePaths=/home/kyle/.pi/agent
InaccessiblePaths=/home/kyle/.ssh /home/kyle/.gnupg /home/kyle/.config/agent
